# jira-comment-slack-server.conf
description "JIRA Comments for Slack"
author  "willharris/SmartXWorks/Kisung"

start on runlevel [234]
stop on runlevel [0156]

setuid jira

exec /usr/local/bin/jira-comment-slack-server

respawn
